# contrib/paavo-web.service
[Unit]
Description=paavo read-only web viewer
After=paavod.service
PartOf=paavod.service

[Service]
Type=simple
User=paavo
Group=paavo
ExecStart=/usr/local/bin/paavo-web --config /etc/paavo/paavo.toml
Restart=on-failure
RestartSec=5
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
